# Round 2 Technical Assessment/Section 1 - Data Pipelines/crontab
0 1 * * * cd / && /usr/local/bin/python pipeline.py >> data/cronlogs.txt
